# deploy/monocle-backend.service
[Unit]
Description=Monocle Backend API
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/monocle/app/agent-backend
EnvironmentFile=/opt/monocle/app/deploy/backend.env
ExecStart=/usr/bin/node dist/app.js
Restart=always
RestartSec=5
User=ubuntu
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monocle-backend

[Install]
WantedBy=multi-user.target
